# menubar/com.claudetracker.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claudetracker.menubar</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexanderdu/Documents/Projects/UsageTracker/menubar/.venv/bin/python</string>
    <string>/Users/alexanderdu/Documents/Projects/UsageTracker/menubar/app.py</string>
  </array>

  <!-- Start at login and restart automatically if it ever crashes -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/claudetracker.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claudetracker.err</string>
</dict>
</plist>
